# deploy/omni-hub.service
[Unit]
Description=Omni-Hub Multi-Tenant WhatsApp Agent Service
Documentation=https://github.com/your-org/omni-hub
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omnihub
Group=omnihub
WorkingDirectory=/opt/omni-hub
Environment="PATH=/opt/omni-hub/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/omni-hub"
EnvironmentFile=/opt/omni-hub/.env

# Main service command
ExecStart=/opt/omni-hub/.venv/bin/uvicorn src.api.app:app --host 0.0.0.0 --port 8000 --workers 1

# Pre-start commands
ExecStartPre=/bin/bash -c 'mkdir -p /opt/omni-hub/logs /opt/omni-hub/data'
ExecStartPre=/opt/omni-hub/.venv/bin/python -c "from src.db.bootstrap import bootstrap_default_instance; bootstrap_default_instance()"

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=0

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/omni-hub/logs /opt/omni-hub/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=omni-hub

[Install]
WantedBy=multi-user.target